# couchbase4/couchbased.conf
[program:couchbase]
command=/opt/couchbase/etc/couchbase_init.d restart
startsecs=5
stopsignal=INT
autostart=true
user=couchbase
redirect_stderr=true

[group:dev]
programs=couchbase