# cassandra3/cassandrad.conf
[program:cassandra]
command=/usr/sbin/cassandra start -f
startsecs=5
stopsignal=INT
autostart=true
user=cassandra
redirect_stderr=true

[group:dev]
programs=cassandra