# tomcat8/tomcatd.conf
[program:tomcat]
command=/opt/tomcat/bin/catalina.sh run
startsecs=5
stopsignal=INT
autostart=true
user=tomcat
redirect_stderr=true

[group:dev]
programs=tomcat