# mysql5/mysqld.conf
[program:mysql]
command=/usr/sbin/mysqld
startsecs=5
stopsignal=INT
autostart=true
user=mysql
redirect_stderr=true

[group:dev]
programs=mysql